# systemd/atom-updater.timer
[Unit]
Description=Atom Updater

[Timer]
Persistent=true
OnUnitActiveSec=10s
OnBootSec=10s

[Install]
WantedBy=timers.target
