# systemd/atom-updater.service
[Unit]
Description=Atom Updater

[Service]
Type=oneshot
ExecStart=/bin/sh -c "atom-updater"

[Install]
WantedBy=multi-user.target
