# systemd/atom-updater.timer
[Unit]
Description=Atom Updater

[Timer]
Persistent=true
OnUnitActiveSec=30min
OnBootSec=30s

[Install]
WantedBy=timers.target
